# ansible/roles/security_audit/templates/lynis-audit.timer.j2
[Unit]
Description=Daily Lynis Security Audit (03:00)
Documentation=https://cisofy.com/lynis/
Requires=lynis-audit.service

[Timer]
# WHY: Run daily at 03:00, before the 04:00 backup.
# This ensures audit results are captured in the daily backup.
OnCalendar=*-*-* 03:00:00
# WHY: Randomize delay to avoid thundering herd on precise :00
RandomizedDelaySec=15min
# WHY: Persist timer state across reboots
Persistent=true

[Install]
WantedBy=timers.target

